# deployment.service
[Unit]
Description=chrissx Media deployment manager

[Service]
ExecStart=/usr/local/bin/deployment

[Install]
WantedBy=multi-user.target
